# ns-bot.conf
[program:ns-bot]
command=python3 bot.py --verbose --log-file bot.log
directory=/home/eigenein/ns-bot
user=eigenein
